# target/linux/oxnas/base-files/etc/init.d/rm-mtdblock
#!/bin/sh /etc/rc.common

START=10

boot() {
	case $(board_name) in
	cloudengines,pogoplug-pro|\
	cloudengines,pogoplug-series-3)
		rm $(find_mtd_part u-boot)
		rm $(find_mtd_part u-boot-env)
		;;
	esac
}
