# target/linux/bcm47xx/base-files/etc/uci-defaults/09_fix_crc
#!/bin/sh
#
# Copyright (C) 2007 OpenWrt.org
#
#

model="$(cat /tmp/sysinfo/model)"

case "$model" in
"ipTIME N704BCM" | \
"ipTIME N904ns")
	exit
	;;
esac

mtd fixtrx firmware
