# target/linux/ath79/base-files/etc/uci-defaults/05_eth_migration
#!/bin/sh

. /lib/functions.sh
. /lib/functions/migrations.sh

board=$(board_name)

case "$board" in
arduino,yun)
	migrate_eths "eth1=eth0"
	;;
esac

migrations_apply network

exit 0
