# target/linux/mediatek/base-files/etc/uci-defaults/99-net-ps
uci set network.globals.packet_steering=1
uci commit
exit 0
